# extras/supervisor.conf
[supervisord]
nodaemon = true

# [program:indexer]
# stdout_logfile = /dev/stdout
# stdout_logfile_maxbytes = 0
# stderr_logfile = /dev/stdout
# stderr_logfile_maxbytes = 0
# command = sh /support/sphinx_indexer.sh

# [program:sphinx]
# stdout_logfile = /dev/stdout
# stdout_logfile_maxbytes = 0
# stderr_logfile = /dev/stdout
# stderr_logfile_maxbytes = 0
# command = sh /support/sphinx_server.sh

## vsftp is not working yet
## "500 OOPS: chroot" error
# [program:vsftpd]
# command = vsftpd
# autostart = true
# autorestart = true
# stdout_logfile = /dev/stdout
# stdout_logfile_maxbytes = 0
# stderr_logfile = /dev/stdout
# stderr_logfile_maxbytes = 0

[program:cron]
command = cron -f -L 15
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stdout
stderr_logfile_maxbytes = 0

[program:apache]
command = /usr/bin/runner.sh
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stdout
stderr_logfile_maxbytes = 0
